# deploy/hosting-panel-update.service
[Unit]
Description=Hosting Panel auto-update from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/opt/hosting-panel/.env
Environment=ENV_FILE=/opt/hosting-panel/.env
Environment=APP_DIR=/opt/hosting-panel
Environment=APP_USER=hosting-panel
Environment=APP_GROUP=hosting-panel
ExecStart=/usr/local/bin/hosting-panel-update
